# supervisord.conf
[supervisord]
nodaemon=true
 
[program:webserver]
command=log.io-server
 
[program:harvester]
command=log.io-harvester